# server/UNITS/after-archive@.service
[Unit]
Description=Convert FLV to MP4 for %i

[Service]
User=www-data
Group=www-data
Type=oneshot
ExecStart=/usr/local/bin/after-archive.sh %i
