# _REPO_SETUP/SERVER_CONFIG/icecast/audio-icecast@.service
[Unit]
Description=Audio Icecast Stream for %i
After=network.target icecast2.service
Requires=icecast2.service

[Service]
Type=simple
User=www-data
Group=www-data
ExecStart=/usr/local/bin/start-audio-icecast.sh %i
ExecStop=/usr/local/bin/stop-audio-icecast.sh %i
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
