# server/stream.service
# /etc/systemd/system/stream.service
[Unit]
Description=SpokenWord FFmpeg Stream
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/start-stream.sh
# перезапускаем **только при сбоях** (код выхода ≠ 0, сигнал, core‑dump)
Restart=on-failure
RestartSec=3
# на всякий случай убираем флаг, когда службу стопят
ExecStopPost=/usr/bin/rm -f /tmp/streaming-enabled

[Install]
WantedBy=multi-user.target


# Загружаем и включаем:
# systemctl daemon-reload
# systemctl enable stream
